# docker/php/supervisor.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=docker-php-entrypoint php-fpm
autostart=true
autorestart=true

[program:queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work --queue=default --sleep=3 --tries=3
autostart=true
autorestart=true
# user=www-data ←docker-compose.yml内で、user: "1000:1000"としていればこの記載は不要
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/worker.log